# scripts/any-kvm-agent.service
[Unit]
Description=Any-KVM Device Agent
Documentation=https://github.com/any-kvm/any-kvm
After=network.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:99
ExecStartPre=-/usr/bin/pkill -f "Xvfb :99"
ExecStart=/usr/bin/any-kvm-agent-wrapper
ExecStopPost=-/usr/bin/pkill -f "Xvfb :99"
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=any-kvm-agent

# 设备访问权限
SupplementaryGroups=video audio input

[Install]
WantedBy=multi-user.target
